# scripts/refresh_data.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job for the daily CSV refresh on macOS.

  Install:
    cp scripts/refresh_data.plist.example \
       ~/Library/LaunchAgents/com.statarb.refresh.plist
    # edit the two absolute paths below to match your checkout
    launchctl load ~/Library/LaunchAgents/com.statarb.refresh.plist
    launchctl start com.statarb.refresh   # one-off test run

  The schedule below is 06:00 local time, weekdays (Mon-Fri). Adjust the
  StartCalendarInterval array for a different cadence.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.statarb.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/YOUR_USERNAME/Desktop/AI Trading/scripts/refresh_data.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/Desktop/AI Trading</string>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/Desktop/AI Trading/logs/refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/Desktop/AI Trading/logs/refresh.err</string>
  <key>RunAtLoad</key>      <false/>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key> <integer>1</integer>
      <key>Hour</key>    <integer>6</integer>
      <key>Minute</key>  <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key> <integer>2</integer>
      <key>Hour</key>    <integer>6</integer>
      <key>Minute</key>  <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key> <integer>3</integer>
      <key>Hour</key>    <integer>6</integer>
      <key>Minute</key>  <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key> <integer>4</integer>
      <key>Hour</key>    <integer>6</integer>
      <key>Minute</key>  <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key> <integer>5</integer>
      <key>Hour</key>    <integer>6</integer>
      <key>Minute</key>  <integer>0</integer>
    </dict>
  </array>
</dict>
</plist>
